# examples/bls/test-thing.workarounds/rhel9/etc/systemd/system/notify-multiuser.service
[Unit]
After=multi-user.target
Wants=multi-user.target

[Service]
LoadCredential=vmm.notify_socket
ExecStart=/etc/notify-multiuser.py
Type=exec
RemainAfterExit=yes
