# examples/bls/test-thing.workarounds/rhel9/etc/systemd/system/sshd-vsock.socket
[Unit]
Wants=ssh-access.target
Before=ssh-access.target

[Socket]
ListenStream=vsock::22
Accept=yes
